# stow/local.example/.config/zsh/profile.d/20-1password-account.zsh
# Local 1Password account settings (template).
# Copy this file into stow/local/.config/zsh/profile.d/ and set your value.

# export OP_ACCOUNT="YOUR_ACCOUNT_ID"
# export ONE_PASSWORD_ACCOUNT="$OP_ACCOUNT"
